# aw-watcher-window.service
[Unit]
Description=aw server

[Service]
ExecStart=/usr/local/bin/activitywatch/aw-watcher-window
Restart=always

[Install]
WantedBy=default.target
