# aw-watcher-afk.service
[Unit]
Description=aw server

[Service]
ExecStart=/usr/local/bin/activitywatch/aw-watcher-afk
Restart=always

[Install]
WantedBy=default.target
